# config/zram.service
[Unit]
Description=Zram

[Service]
ExecStart=/usr/bin/zram.sh

[Install]
WantedBy=multi-user.target
